# config/supervisor.conf
[program:pg_reloaded]
command=/usr/local/pg_reloaded start --config=/etc/pg_reloaded/pg_reloaded.yml
